# otus_pro/borgbackup/templates/backup.timer.j2
[Unit]
Description=Run borg-backup every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=borg-backup.service

[Install]
WantedBy=timers.target
